# supervisor/weather.conf
[program:weather]
command=
directory=
autostart=true
autorestart=true
